# examples/thethingsnetwork-send-v2-zeromq/service/cbridge.service
[Unit]
Description=Commean (TTN) Bridge
After=syslog.target

[Service]
Type=simple
ExecStart=cbridge

[Install]
WantedBy=multi-user.target
